# config/systemd/mpris-player-info-waybar.service
[Unit]

Requires=mpris-player-info-waybar.service

[Service]
ExecStart=%h/.cargo/bin/mpris-player-info waybar
Type=notify
