# config/systemd/playerctld.service

[Unit]
Description = Playerctl state daemon

[Service]
ExecStart=/usr/bin/playerctld
Type=dbus
BusName=org.mpris.MediaPlayer2.playerctld
